# sv/protonvpn-api-server/run
#!/bin/sh
# ABOUTME: ProtonVPN API Server - runit service script
# ABOUTME: HTTP/WebSocket API server with runit supervision and authentication

set -e

# Service configuration
SERVICE_USER="protonvpn"
SERVICE_GROUP="protonvpn"
API_SERVER="/usr/local/bin/api-server"
CONFIG_FILE="/etc/protonvpn/api-server.conf"

# Check if service user exists
if ! getent passwd "$SERVICE_USER" >/dev/null 2>&1; then
    echo "ERROR: Service user '$SERVICE_USER' does not exist"
    exit 1
fi

# Check if API server exists
if [ ! -x "$API_SERVER" ]; then
    echo "ERROR: API server not found or not executable: $API_SERVER"
    exit 1
fi

# Set up environment
export VPN_LOG_DIR="/var/log/protonvpn"
export VPN_RUN_DIR="/run/protonvpn"
export VPN_CONFIG_DIR="/etc/protonvpn"

# Ensure directories exist
[ -d "$VPN_LOG_DIR" ] || mkdir -p "$VPN_LOG_DIR"
[ -d "$VPN_RUN_DIR" ] || mkdir -p "$VPN_RUN_DIR"
chown "$SERVICE_USER:$SERVICE_GROUP" "$VPN_LOG_DIR" "$VPN_RUN_DIR" 2>/dev/null || true

# Log service start
echo "$(date '+%Y-%m-%d %H:%M:%S') [INFO] Starting ProtonVPN API server under runit supervision"

# Execute API server as service user
exec chpst -u "$SERVICE_USER:$SERVICE_GROUP" \
     -e /etc/protonvpn \
     "$API_SERVER" --config "$CONFIG_FILE" --daemon 2>&1
