# systemd/protonvpn-daemon.service
[Unit]
Description=ProtonVPN Configuration Updater Daemon
After=network-online.target
Wants=network-online.target
Documentation=man:protonvpn-daemon(8)

[Service]
Type=forking
User=protonvpn
Group=protonvpn
ExecStart=/usr/local/bin/protonvpn-updater-daemon-secure.sh
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/protonvpn/updater.pid
Restart=on-failure
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=10

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/protonvpn /run/protonvpn /etc/protonvpn
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
PrivateDevices=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
