# systemd/protonvpn-notification.service
[Unit]
Description=ProtonVPN Notification Manager
After=network-online.target
Wants=network-online.target
Documentation=man:notification-manager(8)

[Service]
Type=simple
User=protonvpn
Group=protonvpn
ExecStart=/usr/local/bin/notification-manager --daemon --config /etc/protonvpn/notification.conf
Restart=on-failure
RestartSec=5
TimeoutStartSec=10
TimeoutStopSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/protonvpn /run/protonvpn
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
PrivateDevices=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
