# service/systemd/protonvpn-updater.service
[Unit]
Description=ProtonVPN Config Updater
Documentation=https://github.com/maxrantil/protonvpn-manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nobody
Group=nobody
ExecStart=/home/mqx/workspace/claude-code/vpn/service/systemd/protonvpn-updater-daemon.sh
Restart=always
RestartSec=30

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/protonvpn /run /home/mqx/workspace/claude-code/vpn/locations
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

# Resource limits
LimitNOFILE=1024
MemoryLimit=50M
CPUQuota=10%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=protonvpn-updater

[Install]
WantedBy=multi-user.target
