# sv/protonvpn-notification/run
#!/bin/sh
# ABOUTME: ProtonVPN Notification Manager - runit service script
# ABOUTME: WCAG 2.1 compliant notification system with runit supervision

set -e

# Service configuration
SERVICE_USER="protonvpn"
SERVICE_GROUP="protonvpn"
NOTIFICATION_MANAGER="/usr/local/bin/notification-manager"
CONFIG_FILE="/etc/protonvpn/notification.conf"

# Check if service user exists
if ! getent passwd "$SERVICE_USER" >/dev/null 2>&1; then
    echo "ERROR: Service user '$SERVICE_USER' does not exist"
    exit 1
fi

# Check if notification manager exists
if [ ! -x "$NOTIFICATION_MANAGER" ]; then
    echo "ERROR: Notification manager not found or not executable: $NOTIFICATION_MANAGER"
    exit 1
fi

# Set up environment
export VPN_LOG_DIR="/var/log/protonvpn"
export VPN_RUN_DIR="/run/protonvpn"
export VPN_CONFIG_DIR="/etc/protonvpn"

# Ensure directories exist
[ -d "$VPN_LOG_DIR" ] || mkdir -p "$VPN_LOG_DIR"
[ -d "$VPN_RUN_DIR" ] || mkdir -p "$VPN_RUN_DIR"
chown "$SERVICE_USER:$SERVICE_GROUP" "$VPN_LOG_DIR" "$VPN_RUN_DIR" 2>/dev/null || true

# Log service start
echo "$(date '+%Y-%m-%d %H:%M:%S') [INFO] Starting ProtonVPN notification manager under runit supervision"

# Execute notification manager as service user
exec chpst -u "$SERVICE_USER:$SERVICE_GROUP" \
     -e /etc/protonvpn \
     "$NOTIFICATION_MANAGER" --daemon --config "$CONFIG_FILE" 2>&1
